# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
loglevel=info

[program:web]
command=gunicorn -w 2 -b 0.0.0.0:5000 'sale_monitor.web.app:create_app()'
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true

[program:monitor]
command=python -m sale_monitor.cli.main --products-csv /app/data/products.csv --state-file /app/data/state.json --every %(ENV_CHECK_INTERVAL)s
environment=PYTHONPATH="/app/src"
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true
